# indicator.service
[Unit]
Description=Wifi

[Service]
Type=simple
ExecStart=/usr/local/bin/vindicator /tmp/signal "wifistart.sh"

[Install]
WantedBy=multi-user.target
