# indicator.service
[Unit]
Description=Wifi

[Service]
Type=simple
Restart=always
ExecStart=vindicator /tmp/signal "wifistart.sh"
